# files/etc/uci-defaults/99-argon-dark
#!/bin/sh

uci -q get argon.@global[0] >/dev/null || uci -q add argon global >/dev/null
uci -q set argon.@global[0].mode='dark'
uci -q commit argon

exit 0
